# agent/build_resources/listener_init
#!/bin/bash

### BEGIN INIT INFO
# Provides:          ncpa_listener
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       Enables the NCPA listener service.
### END INIT INFO

_BASEDIR_

# Grab the PID file location
# Someday the py file should get status so this doesn't need to happen
PIDFILE=$(grep pidfile -r $BASEDIR/etc/*.cfg | grep listener | sed '/^#/ d' | cut -d '=' -f 2 | awk '{$1=$1};1' | head -1)
PIDFILE=$BASEDIR/$PIDFILE

# Start listener service
start() {
    (
        cd $BASEDIR
        ./ncpa_posix_listener --start
    )
    echo "Started listener..."
}

# Restart listener service
stop() {
    (
        cd $BASEDIR
        ./ncpa_posix_listener --stop
        echo "Exited."
    )
}

# Get the status
status() {
    if [ -f "$PIDFILE" ];
    then
        pid=`cat $PIDFILE`
        echo "Process running with PID $pid"
    else
        echo "Process is not running."
    fi
}

# Main init script logic
case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status
        ;;
    restart|reload|condrestart)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac
exit 0
