# agent/build_resources/ncpa_listener.service
[Unit]
Description=NCPA Listener
Documentation=https://www.nagios.org/ncpa
After=network.target local-fs.target

[Service]
ExecStart=/usr/local/ncpa/ncpa_listener -n

[Install]
WantedBy=multi-user.target
