# agent/build_resources/passive_init
#!/bin/bash

### BEGIN INIT INFO
# Provides:          ncpa_passive 
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: This manages the NCPA Passive service
# Description:       NCPA is a cross platform monitoring agent for Nagios that requires no dependencies.
### END INIT INFO

# init.d / servicectl compatibility (openSUSE)
if [ -f /etc/rc.status ]; then
    . /etc/rc.status
    rc_reset
fi

# Source function library
if [ -f /etc/rc.d/init.d/functions ]; then
    . /etc/rc.d/init.d/functions
fi

_BASEDIR_

# Grab the PID file location
PIDFILE=$(grep pidfile -r $BASEDIR/etc/*.cfg | grep passive | sed '/^#/ d' | cut -d '=' -f 2 | awk '{$1=$1};1' | head -1)
PIDFILE=$BASEDIR/$PIDFILE

# Start NCPA Passive service
start() {
    PID=`pidofproc -p $PIDFILE`
    if [ checkpid $PID ]; then
        (
            cd $BASEDIR
            action "Starting ncpa_passive" ./ncpa_posix_passive --start
        )
    else
        failure && echo "ncpa_passive is already running ($PID)"
    fi
}

# Stop NCPA Passive service
stop() {
    (
        cd $BASEDIR
        action "Stopping ncpa_passive" ncpa_posix_passive --stop
    )
}

# Get the status of the NCPA Passive service
status() {
    status -p $PIDFILE ncpa_passive
}

# Main init script logic
case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status
        ;;
    restart|reload|condrestart)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac
exit 0
